# conf/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

[program:forever]
command=forever start /usr/src/app/conf/forever/dev.json
environment=NODE_ENV=dev